# sigma_rules/security_medusa_service_creation.yml
title: MedusaLocker Ransomware Scheduled Task Persistence
id: c0f0ad54-e4ef-4456-aa2c-9ba8384d43f7
status: experimental
description: Detects scheduled task creation events (Event ID 4698) that reference the svhost.exe or svhostt.exe executables in the user's Roaming AppData directory, which are characteristic of MedusaLocker ransomware persistence mechanisms.
references:
    - https://www.cisa.gov/news-events/cybersecurity-advisories/aa22-181a
author: Tyler Casey @1qazCasey, @scythe_io
date: 2023/12/29
tags:
    - attack.persistence
    - attack.t1053.005
    - attack.execution
    - attack.t1059
logsource:
    product: windows
    service: security
detection:
    selection:
        EventID: 4698
        TaskContent|contains:
            - 'svhost.exe'
            - 'svhostt.exe'
            - '\AppData\Roaming\'
    condition: selection
falsepositives:
    - Unlikely; the svhost/svhostt naming in AppData is specific to MedusaLocker.
level: high
